# dist/mastodon-streaming.service
[Unit]
Description=tucano-streaming
After=network.target
Wants=tucano-streaming@4000.service

[Service]
Type=oneshot
ExecStart=/bin/echo "tucano-streaming exists only to collectively start and stop tucano-streaming@ instances, shimming over the migration to templated tucano-streaming systemd unit"
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
